# config/cron/fundry.payment.cron
# runs a daily payment processing job for completed features.
10 7 * * * cd /var/www/fundry && ./bin/fundry-payment-features-daily
