# config/cron/sphinx.cron
# add this cronjob on production server
*/3 * * * * sudo indexer --all --rotate > /dev/null
