# config/cron/status.cron
*/10 * * * * cd /var/www/fundry && ./bin/fundry-status > public/__status__.txt
